# misc/autorun/init/overcli-monitor.conf
description     'OverC Monitor'
author  "kolypto <kolypto@gmail.com>"

start on runlevel [2345]
stop on runlevel [016]

chdir /etc/overc

respawn

instance $MONITOR
env MONITOR=

pre-start script
    # No instance name provided: launch multiple
    if [ -z "$MONITOR" ] ; then
        # Scan monitor*.ini files
        for f in `ls /etc/overc/monitor*.ini` ; do
            start overcli-monitor MONITOR=$f
        done
    fi
end script

script
    if [ -z "$MONITOR" ] ; then
        # Fake main process :)
        while true ; do sleep 60 ; done
    else
        # Instance name provided: launch with config file
        overcli monitor $MONITOR
    fi
end script

post-stop script
    if [ -z "$MONITOR" ] ; then
        for f in `ls /etc/overc/monitor*.ini` ; do
            stop overcli-monitor MONITOR=$f
        done
    fi
end script
