# misc/docker/overc-server/conf/supervisor-all.conf
[supervisord]
logfile = /var/log/supervisord.log
logfile_backups = 0

[program:nginx]
command = /usr/sbin/nginx

[program:uwsgi]
command = /usr/local/bin/uwsgi /etc/uwsgi/overc.yml
